# midi-mapper.service
[Unit]
Description=MIDI Mapper Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/midi-mapper-linux --daemon
Restart=always
RestartSec=5

# Proper shutdown
KillSignal=SIGTERM
TimeoutStopSec=10

# Performance tuning - lower latency
Nice=-5

[Install]
WantedBy=default.target